# docker/supervisord-services.ini
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /home/app/django/docker/cla_backend.ini --die-on-term

[program:celery]
command=celery -A cla_backend worker -l info -c 1 --logfile=/var/log/celery/celery.log --pidfile=/var/run/celery/celery.pid
directory=/home/app/django
environment=DJANGO_SETTINGS_MODULE="cla_backend.settings"
