# docker/supervisord-services.ini
[program:nginx]
command=/usr/sbin/nginx -c /home/app/django/docker/nginx.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /home/app/django/docker/cla_backend.ini --die-on-term
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes = 0

# [program:celery]
# command=celery -A cla_backend worker -l info -c 1 --logfile=/dev/stdout --pidfile=/var/run/celery/celery.pid
# directory=/home/app/django
# environment=DJANGO_SETTINGS_MODULE="cla_backend.settings"
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# redirect_stderr=true
