# docker/supervisord-services.ini
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:uwsgi]
command=/usr/sbin/uwsgi --ini /home/app/django/docker/cla_backend.ini --die-on-term --log-master
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:celery]
command=celery -A cla_backend worker -l info -c 1 --logfile=/dev/stdout --pidfile=/var/run/celery/celery.pid
directory=/home/app/django
environment=DJANGO_SETTINGS_MODULE="cla_backend.settings"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
